# B-DOP-400-REN-4-1-octopus-dorian.ayoul/roles/result/files/result.service
[Unit]
Description=Uploads result service

[Service]
User=root
WorkingDirectory=/root/result/
Restart=always
ExecStart=npm start
Environment=DB_NAME="{{ DB_NAME }}" USER_NAME="{{ USER_NAME }}" USER_PASSWORD="{{ USER_PASSWORD }}"